# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:backend]
command=mvn spring-boot:run -f /app/backend/pom.xml

[program:frontend]
command=npm --prefix /app/frontend start /app/frontend
